# systemd/centrix-dashboard.service
[Unit]
Description=Centrix Dashboard API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/centrix
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/centrix/.venv/bin/python -m centrix.dashboard.server
Restart=on-failure

[Install]
WantedBy=default.target
