# systemd/centrix-tui.service
[Unit]
Description=Centrix Control TUI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/centrix
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/centrix/.venv/bin/python -m centrix.tui.control
Restart=on-failure

[Install]
WantedBy=default.target
